# back-end/docker/conf/supervisor.conf
[supervisord]
nodaemon=true

[program:websockets]
command=php artisan websockets:serve
directory=/var/www
numprocs=1
autostart=true
autorestart=true
user=root
