# back-end/websockets.conf
[program:websockets]
command=php artisan websockets:serve
directory=/var/www/websocket-server
numprocs=1
autostart=true
autorestart=true
user=root
